# systemd/dashboard.service
[Unit] 
Description=minecraft dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/srv/dashboard
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
